# elgato.service
[Service]
ExecStart=node /home/rash/automation/index.js
Restart=always
